# deploy/tresamigos-api.service
[Unit]
Description=Tres Amigos API
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/var/www/tresamigos
Environment=NODE_ENV=production
EnvironmentFile=/var/www/tresamigos/.env
ExecStart=/usr/bin/pnpm --filter @tresamigos/api start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
